# conf/etc/init/datomic.conf
start on runlevel [2345]
stop on runlevel [!2345]

env DATOMIC_HOME="/opt/datomic-free-0.9.4899"

exec $DATOMIC_HOME/bin/transactor $DATOMIC_HOME/config/datomic.properties